# config/runit/spice-vdagent/run
#!/bin/sh
exec /usr/bin/spice-vdagentd -x
